# dashboard.conf
description "Dashboard.js"
author "ntim"

start on (runlevel [2345])
stop on (runlevel [!2345])

respawn
respawn limit 3 30

script
exec /bin/bash << EOT
	cd /home/pi/dashboard.js
	PORT=3001 node bin/www
EOT
end script
